# shared/launchd/com.spratt.trip-sync.plist.example
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.spratt.trip-sync</string>
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>-c</string>
        <string>source /Users/YOUR_USER/.config/spratt/infrastructure/env.sh &amp;&amp; exec /usr/bin/python3 /Users/YOUR_USER/.config/spratt/infrastructure/trips/trip-sync.py --scan</string>
    </array>
    <key>WatchPaths</key>
    <array>
        <string>/Users/YOUR_USER/.config/spratt/memory/trips</string>
    </array>
    <key>ThrottleInterval</key>
    <integer>5</integer>
    <key>StandardOutPath</key>
    <string>/Users/YOUR_USER/Library/Logs/spratt/trip-sync.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/YOUR_USER/Library/Logs/spratt/trip-sync.log</string>
</dict>
</plist>
