# shared/launchd/com.spratt.outbox-sender.plist.example
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.spratt.outbox-sender</string>
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>-c</string>
        <string>source /Users/YOUR_USER/.config/spratt/infrastructure/env.sh &amp;&amp; exec /usr/bin/python3 /Users/YOUR_USER/.config/spratt/infrastructure/outbox/sender.py</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/Users/YOUR_USER/Library/Logs/spratt/outbox-sender.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/YOUR_USER/Library/Logs/spratt/outbox-sender.log</string>
    <key>WorkingDirectory</key>
    <string>/Users/YOUR_USER/.config/spratt/infrastructure/outbox</string>
</dict>
</plist>
